# deployment/hermesbridge-postgres-backup.service
[Unit]
Description=Create PostgreSQL backup for HermesBridge
After=postgresql.service
Wants=postgresql.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/hermesbridge-postgres-backup.sh
